# mySite.conf
[program:mySite]
command=/usr/local/bin/gunicorn mySite.wsgi:application -b 0.0.0.0:8000
directory=/var/www/mySite
autostart=true
autorestart=true